# configs/opsen-client.service
[Unit]
Description=Opsen Load Balancer Client
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/opsen-client -config /etc/opsen/client.yml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Create config directory if it doesn't exist
ExecStartPre=/bin/mkdir -p /etc/opsen

[Install]
WantedBy=multi-user.target
